# deploy/launchd/com.hivemind.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for the Hivemind bridge.
  install-service.sh substitutes @ROOT@ with the project absolute path.
  KeepAlive => relaunch on crash; RunAtLoad => start at login.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hivemind.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>@ROOT@/.venv/bin/python</string>
    <string>-m</string>
    <string>hivemind</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@ROOT@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>@ROOT@/var/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>@ROOT@/var/logs/launchd.err.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
